# lib/generators/services/templates/delayed_job.service.erb
[Unit]
Description=<%= @airport_code %> DelayedJob service

[Service]
Type=simple
Restart=always
RestartSec=5s
Environment=RAILS_ENV=production
WorkingDirectory=<%= Rails.root %>
ExecStart=/usr/local/rvm/wrappers/ruby-2.1.1/rake jobs:work

[Install]
WantedBy=multi-user.target
